# scripts/alnsd/skel/etc/local.d/pacman-init.start

#!/bin/bash

echo "Set up pacman keys"
pacman-key --init
pacman-key --populate archlinux
#pacman-key --add /etc/publickey.asc
#pacman-key --lsign-key 97BEEEC2
pacman-key --populate openrc
